# janusgraph-dist/janusgraph-dist-hadoop-2/supervisord.conf
[supervisord]
nodaemon=true

[program:janusgraph]
priority=1
command=./bin/janusgraph.sh -v start

[program:gremlin]
priority=999
command=./bin/gremlin-server.sh
